# supervisor.conf
[program:myproject]
command= sh /src/start.sh
stdout_logfile=/logs/supervisor.log
redirect_stderr=true
autostart=true
autorestart=true